# etc/systemd/check-ip@.timer
[Unit]
Description=check-ip [%i]

[Timer]
Unit=check-ip@%i.service
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
