# etc/systemd/check-remote@.timer
[Unit]
Description=check-remote [%i]
Requires=check-remote@%i.service

[Timer]
OnCalendar=*:*
Persistent=true

[Install]
WantedBy=timers.target
