# etc/systemd/check-ip@.timer
[Unit]
Description=check-ip [%i]
Requires=check-ip@%i.service

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
