# etc/systemd/check-remote@.service
[Unit]
Description=check-remote [%i]
OnFailure=notify-status@%n

[Service]
Type=oneshot
ExecStart=/opt/msladek/stuff/bin/jobs/check-remote.sh %i
