# etc/systemd/check-reboot.timer
[Unit]
Description=check-reboot timer
 
[Timer]
Unit=check-reboot.service
OnCalendar=21:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
