# etc/systemd/fstrim.service
[Unit]
Description=fstrim

[Service]
Type=oneshot
ExecStart=/opt/msladek/stuff/bin/jobs/fstrim.sh
OnFailure=notify-failure@%n
