# etc/systemd/fstrim.timer
[Unit]
Description=fstrim timer
 
[Timer]
Unit=fstrim.service
OnCalendar=Sun 05:35 UTC
Persistent=true

[Install]
WantedBy=timers.target
