# etc/systemd/fstrim.timer
[Unit]
Description=fstrim timer
Requires=fstrim.service
 
[Timer]
OnCalendar=Sun 05:35 UTC
Persistent=true

[Install]
WantedBy=timers.target
