# etc/systemd/hosts-update.service
[Unit]
Description=hosts-update
OnFailure=notify-status@%n

[Service]
Type=oneshot
ExecStart=/opt/msladek/stuff/bin/jobs/hosts-update.sh
