# etc/systemd/hosts-update.service
[Unit]
Description=hosts-update

[Service]
Type=oneshot
ExecStart=/opt/msladek/stuff/bin/jobs/hosts-update.sh
OnFailure=notify-failure@%n
