# etc/systemd/mysql-dump.service
[Unit]
Description=mysql-dump
OnFailure=notify-status@%n

[Service]
Type=oneshot
ExecStart=/opt/msladek/stuff/bin/jobs/mysql-dump.sh
