# etc/systemd/mysql-dump.service
[Unit]
Description=mysql-dump

[Service]
Type=oneshot
ExecStart=/opt/msladek/stuff/bin/jobs/mysql-dump.sh
OnFailure=notify-failure@%n
