# etc/systemd/mysql-dump.timer
[Unit]
Description=mysql-dump timer
Requires=mysql-dump.service
 
[Timer]
OnCalendar=01:15 UTC
Persistent=true

[Install]
WantedBy=timers.target
