# etc/systemd/mysql-dump.timer
[Unit]
Description=mysql-dump timer
 
[Timer]
Unit=mysql-dump.service
OnCalendar=01:15 UTC
Persistent=true

[Install]
WantedBy=timers.target
