# etc/systemd/os-rsync.timer
[Unit]
Description=os-rsync timer
Requires=os-rsync.service
 
[Timer]
OnCalendar=Sun 05:45 UTC
Persistent=true

[Install]
WantedBy=timers.target
