# etc/systemd/os-rsync.service
[Unit]
Description=os-rsync
OnFailure=notify-status@%n

[Service]
Type=oneshot
ExecStart=/opt/msladek/stuff/bin/jobs/os-rsync.sh
