# etc/systemd/smart-dump.service
[Unit]
Description=smart-dump

[Service]
Type=oneshot
ExecStart=/opt/msladek/stuff/bin/jobs/smart-dump.sh
OnFailure=notify-failure@%n
