# etc/systemd/smart-dump.timer
[Unit]
Description=smart-dump timer
Requires=smart-dump.service
 
[Timer]
OnCalendar=Sun 05:40 UTC
Persistent=true

[Install]
WantedBy=timers.target
