# etc/systemd/smart-dump.service
[Unit]
Description=Dumping smart attributes

[Service]
Type=oneshot
ExecStart=/opt/msladek/stuff/bin/jobs/smart-dump.sh
