# etc/systemd/smart-dump.timer
[Unit]
Description=smart-dump timer
 
[Timer]
Unit=smart-dump.service
OnCalendar=Sun 05:40 UTC
Persistent=true

[Install]
WantedBy=timers.target
