# etc/systemd/zfs-health.service
[Unit]
Description=zfs-health

[Service]
Type=oneshot
ExecStart=/opt/msladek/stuff/bin/jobs/zfs-health.sh
