# etc/systemd/zfs-keystatus@.timer
[Unit]
Description=zfs-keystatus [%i]

[Timer]
Unit=zfs-keystatus@%i.service
OnStartupSec=1min

[Install]
WantedBy=timers.target
