# etc/systemd/zfs-health.timer
[Unit]
Description=zfs-health timer
 
[Timer]
Unit=zfs-health.service
OnCalendar=*:30
Persistent=true

[Install]
WantedBy=timers.target
