# etc/systemd/zfs-health.service
[Unit]
Description=zfs-health
OnFailure=notify-status@%n

[Service]
Type=oneshot
ExecStart=/opt/msladek/stuff/bin/jobs/zfs-health.sh
