# etc/systemd/zfs-scrub.timer
[Unit]
Description=zfs-scrub timer
 
[Timer]
Unit=zfs-scrub.service
OnCalendar=*-*-7,22 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
